# scripts/launchd/com.shawnos.scout.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.shawnos.scout</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/python3</string>
        <string>/Users/shawnos.ai/shawn-gtme-os/scripts/agents/scout.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/shawnos.ai/shawn-gtme-os</string>

    <key>RunAtLoad</key>
    <false/>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>6</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/shawnos.ai</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/shawnos.ai/.logs/scout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/shawnos.ai/.logs/scout.log</string>
</dict>
</plist>
