# scripts/launchd/com.shawnos.security-audit.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.shawnos.security-audit</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/shawnos.ai/shawn-gtme-os/scripts/security_audit.sh</string>
  </array>

  <!-- Every 72 hours: Mon 2am, Thu 2am, Sun 2am -->
  <key>StartCalendarInterval</key>
  <array>
    <dict>
      <key>Weekday</key>
      <integer>1</integer>
      <key>Hour</key>
      <integer>2</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>
    <dict>
      <key>Weekday</key>
      <integer>4</integer>
      <key>Hour</key>
      <integer>2</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>
    <dict>
      <key>Weekday</key>
      <integer>0</integer>
      <key>Hour</key>
      <integer>2</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/shawnos.ai/shawn-gtme-os</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/Users/shawnos.ai/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>HOME</key>
    <string>/Users/shawnos.ai</string>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/shawnos.ai/shawn-gtme-os/data/daily-log/security-audit.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/shawnos.ai/shawn-gtme-os/data/daily-log/security-audit-error.log</string>

</dict>
</plist>
